# sdfs-dd-mgr.service
[Unit]
Description=sdfs dd-mgr

[Service]
Type=notify
ExecStart=/usr/sbin/sdfs-dd-mgr
